# ChickCtrlApp.service
[Unit]
Description=ChickApp Unit

[Service]
ExecStart=/usr/local/bin/chickApp/ChickenCtrlMainApp

[Install]
WantedBy=multi-user.target
